# config/com.dugout.fpl-pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
macOS LaunchAgent for automated FPL data pulls.

Installation:
1. Copy this file to ~/Library/LaunchAgents/
2. Update the paths below to match your setup
3. Load with: launchctl load ~/Library/LaunchAgents/com.dugout.fpl-pull.plist
4. Unload with: launchctl unload ~/Library/LaunchAgents/com.dugout.fpl-pull.plist

This runs every 6 hours. The script itself checks if 5 hours have passed since GW deadline.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dugout.fpl-pull</string>
    
    <key>ProgramArguments</key>
    <array>
        <!-- Update this path to your Python -->
        <string>/Users/YOUR_USERNAME/.pyenv/versions/3.10.14/bin/python</string>
        <string>/Users/YOUR_USERNAME/path/to/the-dugout/scripts/pull_fpl_data.py</string>
        <string>--db</string>
        <string>/Users/YOUR_USERNAME/path/to/the-dugout/storage/fpl_2025_26.sqlite</string>
    </array>
    
    <!-- Run every 6 hours -->
    <key>StartInterval</key>
    <integer>21600</integer>
    
    <!-- Also run at startup -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/path/to/the-dugout</string>
    
    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/YOUR_USERNAME/path/to/the-dugout/src</string>
    </dict>
    
    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/dugout-fpl-pull.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/dugout-fpl-pull-error.log</string>
</dict>
</plist>
